# deploy/systemd/hermes-investment-scheduler.service
[Unit]
Description=Hermes investment assistant scheduler daemon
After=network-online.target investment-assistant-postgres.service
Wants=network-online.target investment-assistant-postgres.service

[Service]
Type=simple
User=jianjustin
Group=jianjustin
WorkingDirectory=/opt/hermes-investment-assistant/app
EnvironmentFile=/opt/hermes-investment-assistant/.env
ExecStart=/opt/hermes-investment-assistant/.venv/bin/python -m investment_assistant.tasks.scheduler
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
