# deploy/systemd/hermes-investment-daily.timer
[Unit]
Description=Daily Hermes investment assistant timer

[Timer]
OnCalendar=Mon..Fri *-*-* 08:30:00 America/New_York
Persistent=true
Unit=hermes-investment-daily.service

[Install]
WantedBy=timers.target
