# deploy/systemd/hermes-investment-scores.timer
[Unit]
Description=Nightly Hermes strategy score timer

[Timer]
OnCalendar=Mon..Fri *-*-* 18:00:00 America/New_York
Persistent=true
Unit=hermes-investment-scores.service

[Install]
WantedBy=timers.target
